# proxy/llama-proxy.user.service
[Unit]
Description=LLama Proxy Server (user)
Documentation=https://github.com/rgardler/llama
# Ensure networking and container socket are up before starting
After=network-online.target podman.socket
Wants=network-online.target
Requires=podman.socket

[Service]
Type=simple
WorkingDirectory=/home/rgardler/projects/llm/proxy
# Ensure PATH contains your venv and user-local bins
Environment="PATH=/home/rgardler/.local/bin:/usr/local/bin:/usr/bin:/bin:/home/rgardler/projects/llm/proxy/.venv/bin"
Environment="LLAMA_PROXY_CONFIG=/home/rgardler/projects/llm/proxy/config.yaml"

# ROCm/AMD GPU environment variables (if needed)
Environment="HSA_OVERRIDE_GFX_VERSION=11.5.1"
Environment="ROCM_LLVM_PRE_VEGA=1"

ExecStart=/home/rgardler/projects/llm/proxy/.venv/bin/python -m uvicorn server:app --host 0.0.0.0 --port 8000
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10

# Logging to the user journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=llama-proxy

# Minimal hardening (user units already isolate from system)
# Note: Do NOT set NoNewPrivileges here because distrobox/newuidmap requires
# the setuid root helper to be able to write uid_map/gid_map during namespace
# setup. If NoNewPrivileges is set, the setuid bit is ignored and namespace
# creation will fail with "newuidmap: write to uid_map failed: Operation not permitted".

# Allow writing to project/user container caches and logs
ReadWritePaths=/home/rgardler/projects/llm/proxy
ReadWritePaths=/home/rgardler/.local/share/containers
ReadWritePaths=/home/rgardler/.cache
ReadWritePaths=/home/rgardler/.config

[Install]
WantedBy=default.target
